# packaging/debian_amd64/lib/systemd/system/cnb-rates.path
[Unit]
Description=openbank CNB Rates services namespace
PartOf=cnb-rates.service

[Path]
PathChanged=/etc/init/cnb-rates.conf

[Install]
WantedBy=cnb-rates.service
