# packaging/debian_amd64/lib/systemd/system/cnb-rates-rest.path
[Unit]
Description=openbank CNB Rates REST service

[Path]
PathChanged=/etc/init/cnb-rates.conf

[Install]
WantedBy=cnb-rates.service
