# packaging/debian/cnb-rates.path
[Unit]
Description=openbank CNB Rates services namespace
PartOf=cnb-rates.service
ConditionPathExists=/etc/cnb-rates/conf.d/init.conf

[Path]
PathChanged=/etc/cnb-rates/conf.d/init.conf
PathModified=/etc/cnb-rates/conf.d/init.conf
Unit=cnb-rates.service

[Install]
WantedBy=cnb-rates.service
