# packaging/debian/cnb-rates-watcher.path
[Unit]
Description=openbank CNB Rates configuration watcher
After=cnb-rates.service

[Path]
PathChanged=/etc/cnb-rates/conf.d
Unit=cnb-rates-watcher.service

[Install]
WantedBy=paths.target
